# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
; --forwarded-allow-ips='*' trusts the X-Forwarded-Proto header from whatever peer connects to
; the container. Cloud Run terminates TLS at its own front end and always proxies to the
; container over a private, single-hop connection - the container is never reachable except
; through that proxy - so this is safe here. Without it, uvicorn's default (only trust
; 127.0.0.1) ignores X-Forwarded-Proto on Cloud Run's real proxy IP.
;
; Wrapped in /bin/sh -c (unlike the monolith's supervisord.conf, which hardcodes --port 8000)
; because this service listens on Cloud Run's injected $PORT (see Dockerfile), and supervisord's
; command= does not itself perform shell/env-var expansion - it must be handed to a shell for
; ${PORT:-8080} to expand. Mirrors the Dockerfile's own CMD before this slice introduced
; supervisord.
command=/bin/sh -c "/app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port ${PORT:-8080} --proxy-headers --forwarded-allow-ips='*'"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=/app/.venv/bin/celery -A app.worker worker --loglevel=info
; Unlike the monolith (organize-me), where this stays autostart=false pending REDIS_URL being
; deployed, this service has REDIS_URL wired for real via --set-secrets in ci.yml/deploy.yml
; (redis-url-{qa,prod} in GCP Secret Manager), so the worker can start immediately alongside web.
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
